# systemd/.config/systemd/user/gpg-refresh-keys.service
[Unit]
Description=GnuPG Refresh Keys
After=network.target network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/gpg --refresh-keys