# systemd/.config/systemd/user/mbsync.timer
[Unit]
Description=Mailbox synchronization timer

[Timer]
OnCalendar=*:0/2
Persistent=true
Unit=mbsync.service

[Install]
WantedBy=timers.target
