# systemd/.config/systemd/user/mbsync.timer
[Unit]
Description=Mailbox synchronization timer

[Timer]
OnBootSec=5s
OnUnitActiveSec=5m
Unit=mbsync.service

[Install]
WantedBy=timers.target
