# systemd/.config/systemd/user/mbsync.service
[Unit]
Description=Mailbox synchronization service
After=network.target network-online.target dbus.socket

[Service]
Type=oneshot
ExecStart=/usr/bin/mbsync -c /home/someone/.config/isync/mbsyncrc -a
StandardOutput=syslog
StandardError=syslog